# systemd/arcade.service
[Unit]
Description=Homelab Arcade supervisor

[Service]
Type=simple
User=arcade
Group=arcade
StateDirectory=homelab-arcade
WorkingDirectory=/var/lib/homelab-arcade
EnvironmentFile=/etc/arcade/arcade.env
ExecStart=/usr/local/bin/homelab-arcade
Restart=on-failure
RestartSec=5s
KillSignal=SIGTERM
TimeoutStopSec=30s
StandardOutput=journal
StandardError=journal
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
