# systemd/arcade.service
[Unit]
Description=Homelab Arcade supervisor

[Service]
Type=simple
User=arcade
Group=arcade
WorkingDirectory=/srv/arcade
EnvironmentFile=/etc/arcade/arcade.env
ExecStart=/srv/arcade/scripts/up.sh
Restart=on-failure
RestartSec=5s
KillSignal=SIGTERM
TimeoutStopSec=30s
StandardOutput=journal
StandardError=journal
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
